# gta02-dm1/rootfs/target_rootfs_data/etc/inittab
# This is run first except when booting in single-user mode.
#
::sysinit:/etc/init.d/rcS

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/halt
::shutdown:/bin/umount -a

